# alpine/signageos.initd
#!/sbin/openrc-run

name=$RC_SVCNAME
description="signageOS client"
pidfile="/run/$name.pid"
command="node"
command_args="/usr/lib/signageos-display-linux/server/server.js"
command_user="signageos"
stopsig=SIGTERM

depend() {
    need net localmount loopback
    use udev
    after nginx
}
