# etc/crontabs/root
*	*	*	*	*	/sbin/signageos-checkserveralive 2>&1 | /usr/bin/logger --tag signageos-checkserveralive
