# common/init.d/set-ondemand-cpu-governor
#!/sbin/openrc-run
# more info here:
# https://putokaz.wordpress.com/2015/03/01/overclocking-raspberry-pi-set-scaling_governor-on-raspbian-boot/
# https://megakemp.com/2013/02/26/adventures-in-overclocking-a-raspberry-pi/

name=$RC_SVCNAME
description="Set ondemand CPU governor and set it's up threshold to 60%"

start() {
    echo ondemand > /sys/devices/system/cpu/cpufreq/policy0/scaling_governor
    echo 60 > /sys/devices/system/cpu/cpufreq/ondemand/up_threshold
    eend 0
}

stop() {
    eend 0
}
