# common/init.d/settermblank
#!/sbin/openrc-run

name=$RC_SVCNAME
description="Forcefully disable on-screen terminal"

start() {
    setterm --blank force
    eend 0
}

stop() {
    eend 0
}
