# common/init.d/signageos-usbmountall
#!/sbin/openrc-run

name=$RC_SVCNAME
description="Auto-mount all USB storage devices on boot"

depend() {
    need localmount
}

start() {
    ebegin "mount all usb devices"
    /sbin/signageos-usbmountall | /usr/bin/logger --tag $name
    eend $? "failed to mount usb devices"
}

stop() {
    # do nothing, but this has to be there because system will call "stop" on this service when switching runlevels
    eend 0
}
